# OS/debian/overlay/etc/systemd/system/serial-getty@ttyPS0.service.d/override.conf
[Unit]
After=rc-local.service multi-user.target

[Service]
Type=simple
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I 115200,9600 linux
