# OS/debian/overlay/etc/systemd/system/hostapd@.path
[Path]
PathExists   = /opt/redpitaya/hostapd.conf
PathModified = /opt/redpitaya/hostapd.conf
