# OS/debian/overlay/etc/systemd/system/wpa_supplicant@.path
[Path]
PathExists   = /opt/redpitaya/wpa_supplicant.conf
PathModified = /opt/redpitaya/wpa_supplicant.conf
