# OS/debian/overlay/etc/profile.d/alias.sh
# Aliases
alias ll='ls -al'
alias df='df -h'
alias rm='rm -i'
alias mv='mv -i'
alias cp='cp -i'
