# OS/debian/overlay/etc/profile.d/redpitaya.sh
export PATH_REDPITAYA=/opt/redpitaya
export PATH=$PATH:$PATH_REDPITAYA/sbin:$PATH_REDPITAYA/bin
#export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:$PATH_REDPITAYA/lib
